# src/agent/vord-agent.service
[Unit]
Description=Vord Agent - Machine Telemetry Collector
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
ExecStart=/usr/local/bin/vord-agent
EnvironmentFile=-/etc/framlux/vord-agent.env
Restart=on-failure
RestartSec=10
TimeoutStopSec=30
MemoryMax=512M
CPUQuota=50%
StateDirectory=vord-agent

ProtectSystem=strict
ReadWritePaths=/var/lib/vord-agent /etc/framlux
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
